# ext/QVD-UP/bin/system/qvd-up-api.service
[Unit]
Description=QVD UP API service
Wants=postgresql.service
After=postgresql.service syslog.target network.target

[Service]
Type=forking
PIDFile=/var/run/qvd/qvd-up-api.pid
ExecStart=/usr/lib/qvd/bin/hypnotoad /usr/lib/qvd/bin/qvd-up-api.pl

[Install]
WantedBy=default.target
