# ext/QVD-VMKiller/cron.d/qvd-vm-killer
01 * * * * root /usr/lib/qvd/bin/qvd-vm-killer
