# ext/QVD-API/bin/qvd-api.service
[Unit]
Description=QVD API service
Wants=postgresql.service
After=postgresql.service syslog.target network.target

[Service]
Type=forking
PIDFile=/var/run/qvd/qvd-api.pid
ExecStart=/usr/lib/qvd/bin/hypnotoad /usr/lib/qvd/bin/qvd-api.pl

[Install]
WantedBy=default.target
